# ina_wah.service
[Unit]
Description=Ina Wah Service

[Service]
Type=simple
ExecStart=/bin/bash /opt/ina_wah.sh
Restart=always

[Install]
WantedBy=multi-user.target
